# src/plasma-kded.service.in
[Unit]
Description=KDE Daemon
PartOf=graphical-session.target
After=plasma-kcminit.service

[Service]
ExecStart=@KDE_INSTALL_FULL_BINDIR@/kded5
BusName=org.kde.kded5
Slice=session.slice
